# harness/ops/launchd/com.tinycloud.distillery.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.tinycloud.distillery.server — the feed server, kept alive across reboots
  (spec §7b). Runs harness/ops/launchd/server.sh, which sources server.env (PATH +
  OPENKEY_ALLOWED_ADDRESSES + TRANSCRIPT_DIRS for the Generate button) then
  builds web/dist and starts `bun src/server.ts` in feed/. KeepAlive +
  RunAtLoad so a reboot brings the feed back without Hunter. __REPO__ / __HOME__
  filled by install-launchd.sh.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tinycloud.distillery.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__REPO__/harness/ops/launchd/server.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO__/harness/feed</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__HOME__</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <!-- Don't hammer-restart a crash-looping server: 10s between respawns. -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>__REPO__/harness/ops/launchd/logs/server.out.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO__/harness/ops/launchd/logs/server.err.log</string>
</dict>
</plist>
