# harness/ops/launchd/com.tinycloud.distillery.tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.tinycloud.distillery.tunnel — the cloudflared named tunnel, kept alive
  across reboots (spec §7b). Runs the `distillery` tunnel
  (64a191a4-9537-413f-83fc-ab1f374774de) from config ~/.cloudflared/distillery.yml,
  which routes distillery.tinytunnel.xyz → localhost:4242. SEPARATE plist from
  the server (spec §7b decision: two plists) so the tunnel can restart
  independently of the server. KeepAlive + RunAtLoad. __HOME__ / __CLOUDFLARED__
  filled by install-launchd.sh (__CLOUDFLARED__ = the absolute cloudflared path).
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tinycloud.distillery.tunnel</string>

    <key>ProgramArguments</key>
    <array>
        <string>__CLOUDFLARED__</string>
        <string>tunnel</string>
        <string>--config</string>
        <string>__HOME__/.cloudflared/distillery.yml</string>
        <string>run</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>TUNNEL_ORIGIN_CERT</key>
        <string>__HOME__/.cloudflared/cert.pem</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>__REPO__/harness/ops/launchd/logs/tunnel.out.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO__/harness/ops/launchd/logs/tunnel.err.log</string>
</dict>
</plist>
