# harness/ops/launchd/com.tinycloud.distillery.feedrun.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.tinycloud.distillery.feedrun — the weekday-morning heartbeat (spec §7a).
  Runs harness/ops/launchd/feedrun.sh, which sources feedrun.env + .env, then invokes the
  feed-run recipe headlessly via `claude -p`. NOT KeepAlive: it is a scheduled
  one-shot, not a daemon. __REPO__ / __HOME__ are filled by install-launchd.sh.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tinycloud.distillery.feedrun</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__REPO__/harness/ops/launchd/feedrun.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO__</string>

    <!-- Weekday mornings at 07:00. launchd has a minimal PATH/env; the wrapper
         rebuilds everything from feedrun.env, so EnvironmentVariables stays
         deliberately bare (only HOME, which a few tools read for ~/.cloudflared
         etc.). -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__HOME__</string>
    </dict>

    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <!-- Don't fire a backlog of missed runs if the laptop was asleep at 7am —
         one catch-up run at wake is enough. -->
    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__REPO__/harness/ops/launchd/logs/feedrun.out.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO__/harness/ops/launchd/logs/feedrun.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
